# response/New folder/9QoorXJoAl_T1547.001_attackPhase.yaml
1. **Framework Section**  
   Phase 1: *Initial Execution*  
   - Monitor for modifications to the following Registry keys:  
     - `HKEY_CURRENT_USER\Software\Microsoft\Windows\CurrentVersion\Run`  
     - `HKEY_LOCAL_MACHINE\Software\Microsoft\Windows\CurrentVersion\Run`  
   - Detect the creation of new files in the user-specific and system-wide Startup folders:  
     - `%APPDATA%\Microsoft\Windows\Start Menu\Programs\Startup`  
     - `%ProgramData%\Microsoft\Windows\Start Menu\Programs\Startup`  
   - Analyze command-line arguments for processes that may indicate the execution of malicious payloads during user logon.  

   Phase 2: *Persistence Mechanism*  
   - Identify and log changes to Registry keys associated with startup programs, particularly:  
     - `HKEY_CURRENT_USER\Software\Microsoft\Windows\CurrentVersion\RunOnce`  
     - `HKEY_LOCAL_MACHINE\Software\Microsoft\Windows\CurrentVersion\RunOnce`  
   - Monitor for the presence of executable files or scripts in the Startup folders that may be linked to malicious activity.  
   - Track the use of the `RunServices` Registry keys to detect attempts to establish persistence through service execution:  
     - `HKEY_LOCAL_MACHINE\Software\Microsoft\Windows\CurrentVersion\RunServices`  

   Phase 3: *Post-Execution Analysis*  
   - Conduct forensic analysis on detected Registry modifications to identify potential malicious entries.  
   - Review system logs for unusual login patterns or execution of programs that coincide with Registry modifications.  
   - Correlate findings with known threat actor TTPs to assess the potential impact and scope of the persistence mechanism.  

2. **Behavioral Patterns Section**  
   - **Pattern Name**: Registry Run Key Modification  
     - Examples:  
       - "ADVSTORESHELL achieves persistence by adding itself to the HKCU\SOFTWARE\Microsoft\Windows\CurrentVersion\Run Registry key."  
       - "Agent Tesla can add itself to the Registry as a startup program to establish persistence."  
     - Associated Phase: Phase 1  
     - Mapping to TTP ID: T1547.001  

   - **Pattern Name**: Startup Folder Manipulation  
     - Examples:  
       - "Amadey has changed the Startup folder to the one containing its executable by overwriting the registry keys."  
       - "APT28 has deployed malware that has copied itself to the startup directory for persistence."  
     - Associated Phase: Phase 1  
     - Mapping to TTP ID: T1547.001  

   - **Pattern Name**: Use of RunOnce for Persistence  
     - Examples:  
       - "APT41 created and modified startup files for persistence."  
       - "AvosLocker has been executed via the RunOnce Registry key to run itself on safe mode."  
     - Associated Phase: Phase 2  
     - Mapping to TTP ID: T1547.001  

   - **Pattern Name**: Service Execution via Registry  
     - Examples:  
       - "APT32 established persistence using Registry Run keys, both to execute PowerShell and VBS scripts."  
       - "Cobalt Group has used Registry Run keys for persistence."  
     - Associated Phase: Phase 2  
     - Mapping to TTP ID: T1547.001  

   - **Pattern Name**: Shortcut Creation in Startup  
     - Examples:  
       - "DarkGate installation finishes with the creation of a registry Run key."  
       - "Trojan.Karagany can create a link to itself in the Startup folder to automatically start itself upon system restart."  
     - Associated Phase: Phase 1  
     - Mapping to TTP ID: T1547.001  