# response/New folder/qvbfrZv2vI_T1053.005_attackPhase.yaml
1. **Framework Section**  
   Phase 1: *Task Creation and Configuration*  
   - Monitor the use of `schtasks.exe` and PowerShell Cmdlets like `New-ScheduledTask` and `Register-ScheduledTask`.  
   - Analyze the creation of scheduled tasks with unusual names or execution paths, especially those that run under SYSTEM or other high-privilege accounts.  
   - Investigate the use of XML files for task definitions, particularly those that may be obfuscated or contain suspicious commands.  

   Phase 2: *Execution and Persistence*  
   - Track the execution of scheduled tasks, focusing on their triggers (e.g., on logon, daily, or at startup).  
   - Identify patterns of repeated execution, such as tasks set to run every minute or hour, which may indicate malicious intent.  
   - Examine the registry keys associated with scheduled tasks, particularly for modifications that hide tasks or alter their visibility.  

   Phase 3: *Artifact and Anomaly Detection*  
   - Detect anomalies in task metadata, such as unexpected changes to the Security Descriptor (SD) or Index values in registry keys.  
   - Monitor for the presence of hidden tasks that do not appear in standard queries (e.g., `schtasks /query`).  
   - Analyze logs for unusual command-line arguments or execution contexts that deviate from normal operational patterns.  

2. **Behavioral Patterns Section**  
   - **Pattern Name**: Scheduled Task Creation  
     - Examples: "APT3 downloader creates persistence by creating the following scheduled task: schtasks /create /tn 'mysc' /tr C:\\Users\\Public\\test.exe /sc ONLOGON /ru 'System'."  
     - Associated Phase: Task Creation and Configuration  
     - Mapping to TTP ID: T1053.005  

   - **Pattern Name**: High-Frequency Execution  
     - Examples: "BADNEWS creates a scheduled task to establish by executing a malicious payload every subsequent minute."  
     - Associated Phase: Execution and Persistence  
     - Mapping to TTP ID: T1053.005  

   - **Pattern Name**: Use of Hidden Tasks  
     - Examples: "An adversary may hide a task from schtasks /query and the Task Scheduler by deleting the associated Security Descriptor (SD) registry value."  
     - Associated Phase: Artifact and Anomaly Detection  
     - Mapping to TTP ID: T1053.005  

   - **Pattern Name**: Task Metadata Manipulation  
     - Examples: "Adversaries may also employ alternate methods to hide tasks, such as altering the metadata (e.g., Index value) within associated registry keys."  
     - Associated Phase: Artifact and Anomaly Detection  
     - Mapping to TTP ID: T1053.005  

   - **Pattern Name**: Remote Task Scheduling  
     - Examples: "During Operation CuckooBees, the threat actors used scheduled tasks to execute batch scripts for lateral movement."  
     - Associated Phase: Execution and Persistence  
     - Mapping to TTP ID: T1053.005  

   - **Pattern Name**: Use of Legitimate Tools for Malicious Purposes  
     - Examples: "Lazarus Group created scheduled tasks to set a periodic execution of a remote XSL script."  
     - Associated Phase: Task Creation and Configuration  
     - Mapping to TTP ID: T1053.005  

   - **Pattern Name**: Scheduled Task for Payload Execution  
     - Examples: "BlackByte created scheduled tasks for payload execution."  
     - Associated Phase: Execution and Persistence  
     - Mapping to TTP ID: T1053.005  

   - **Pattern Name**: Scheduled Task with Obfuscated Commands  
     - Examples: "Chimera has used scheduled tasks to invoke Cobalt Strike including through batch scripts."  
     - Associated Phase: Task Creation and Configuration  
     - Mapping to TTP ID: T1053.005  