# response/New folder (2)/T1053.005_2fXauJ1j19_attackPhaseAnalysis_COSTAR_attackPhase.yaml
1. **Framework Section**  
   Phase 1: *Task Creation and Configuration*  
   - Monitor the use of `schtasks.exe` and other methods (e.g., WMI, PowerShell Cmdlets) for creating scheduled tasks.  
   - Analyze command-line arguments for task creation, such as `/create`, `/tn`, `/tr`, `/sc`, and `/ru`.  
   - Investigate registry keys related to scheduled tasks, particularly those that may indicate hidden tasks (e.g., Security Descriptor values).  

   Phase 2: *Execution and Persistence*  
   - Track the execution of scheduled tasks, focusing on the context in which they run (e.g., SYSTEM, user accounts).  
   - Identify patterns of recurring execution, such as tasks set to run at startup or on specific schedules (e.g., daily, hourly).  
   - Examine logs for task execution events and correlate them with known malicious activity.  

   Phase 3: *Artifact and Anomaly Detection*  
   - Detect anomalies in task names and configurations that may indicate obfuscation or malicious intent (e.g., names mimicking legitimate tasks).  
   - Monitor for modifications to existing tasks, especially those that may indicate lateral movement or privilege escalation.  
   - Analyze the presence of tasks that are not visible through standard enumeration methods, indicating potential evasion techniques.  

2. **Behavioral Patterns Section**  
   - **Pattern Name**: Scheduled Task Creation  
     - Examples: "APT3 downloader creates persistence by creating the following scheduled task: schtasks /create /tn 'mysc' /tr C:\\Users\\Public\\test.exe /sc ONLOGON /ru 'System'."  
     - Associated Phase: Phase 1  
     - Mapping to TTP ID: T1053.005  

   - **Pattern Name**: Recurring Execution  
     - Examples: "BADNEWS creates a scheduled task to establish by executing a malicious payload every subsequent minute."  
     - Associated Phase: Phase 2  
     - Mapping to TTP ID: T1053.005  

   - **Pattern Name**: Hidden Task Manipulation  
     - Examples: "An adversary may hide a task from schtasks /query and the Task Scheduler by deleting the associated Security Descriptor (SD) registry value."  
     - Associated Phase: Phase 3  
     - Mapping to TTP ID: T1053.005  

   - **Pattern Name**: Task Context Execution  
     - Examples: "Chimera has used scheduled tasks to invoke Cobalt Strike including through batch scripts schtasks /create /ru 'SYSTEM'."  
     - Associated Phase: Phase 2  
     - Mapping to TTP ID: T1053.005  

   - **Pattern Name**: Task Name Obfuscation  
     - Examples: "Maze has created scheduled tasks using name variants such as 'Windows Update Security' to launch Maze at a specific time."  
     - Associated Phase: Phase 3  
     - Mapping to TTP ID: T1053.005  

   - **Pattern Name**: Remote Task Scheduling  
     - Examples: "During Operation CuckooBees, the threat actors used scheduled tasks to execute batch scripts for lateral movement."  
     - Associated Phase: Phase 2  
     - Mapping to TTP ID: T1053.005  

   - **Pattern Name**: Task Modification  
     - Examples: "APT29 manipulated scheduled tasks by updating an existing legitimate task to execute their tools."  
     - Associated Phase: Phase 3  
     - Mapping to TTP ID: T1053.005  