# response/New folder (2)/T1547.001_daIw6X7Pnk_attackPhaseAnalysis_COSTAR_attackPhase.yaml
1. **Framework Section**  
   Phase 1: *Initial Execution*  
   - Monitor for modifications to the following Registry keys:  
     - `HKEY_CURRENT_USER\Software\Microsoft\Windows\CurrentVersion\Run`  
     - `HKEY_LOCAL_MACHINE\Software\Microsoft\Windows\CurrentVersion\Run`  
   - Detect the creation of executable files in the Startup folder:  
     - `%APPDATA%\Microsoft\Windows\Start Menu\Programs\Startup`  
     - `%ProgramData%\Microsoft\Windows\Start Menu\Programs\Startup`  
   - Analyze command-line arguments for processes that may indicate the execution of malicious payloads during user login.  

   Phase 2: *Persistence Mechanism*  
   - Identify and log changes to the Registry keys associated with startup programs:  
     - `HKEY_CURRENT_USER\Software\Microsoft\Windows\CurrentVersion\RunOnce`  
     - `HKEY_LOCAL_MACHINE\Software\Microsoft\Windows\CurrentVersion\RunOnce`  
   - Monitor for the presence of shortcuts (`.lnk` files) in the Startup folder that may point to malicious executables.  
   - Track the use of the `RunServices` Registry keys to detect attempts to establish persistence through service execution:  
     - `HKEY_LOCAL_MACHINE\Software\Microsoft\Windows\CurrentVersion\RunServices`  

   Phase 3: *Malware Execution*  
   - Analyze the execution context of processes initiated from the Startup folder or Registry keys to identify potential privilege escalation or unauthorized access.  
   - Monitor for unusual network activity or connections initiated by processes that have been executed via persistence mechanisms.  
   - Investigate the integrity of files referenced in Registry keys or Startup folder entries for signs of tampering or malicious behavior.  

2. **Behavioral Patterns Section**  
   - **Pattern Name**: Registry Run Key Modification  
     - Examples:  
       - "ADVSTORESHELL achieves persistence by adding itself to the HKCU\SOFTWARE\Microsoft\Windows\CurrentVersion\Run Registry key."  
       - "Agent Tesla can add itself to the Registry as a startup program to establish persistence."  
     - Associated Phase: Phase 1  
     - Mapping to TTP ID: T1547.001  

   - **Pattern Name**: Startup Folder Manipulation  
     - Examples:  
       - "Amadey has changed the Startup folder to the one containing its executable by overwriting the registry keys."  
       - "APT28 has deployed malware that has copied itself to the startup directory for persistence."  
     - Associated Phase: Phase 2  
     - Mapping to TTP ID: T1547.001  

   - **Pattern Name**: Service Execution via Registry  
     - Examples:  
       - "APT32 established persistence using Registry Run keys, both to execute PowerShell and VBS scripts as well as to execute their backdoor directly."  
       - "APT41 created and modified startup files for persistence."  
     - Associated Phase: Phase 3  
     - Mapping to TTP ID: T1547.001  

   - **Pattern Name**: Shortcut Creation in Startup  
     - Examples:  
       - "DarkGate installation finishes with the creation of a registry Run key."  
       - "Trojan.Karagany can create a link to itself in the Startup folder to automatically start itself upon system restart."  
     - Associated Phase: Phase 2  
     - Mapping to TTP ID: T1547.001  

   - **Pattern Name**: Registry Key for Automatic Execution  
     - Examples:  
       - "Cobalt Group has used Registry Run keys for persistence."  
       - "Ursnif has used Registry Run keys to establish automatic execution at system startup."  
     - Associated Phase: Phase 1  
     - Mapping to TTP ID: T1547.001  